# systemd/example-app.service.d/crashomon.conf
# Drop-in snippet for any systemd service that wants crash monitoring.
#
# Copy this file to /etc/systemd/system/<your-service>.service.d/crashomon.conf
# and run: systemctl daemon-reload && systemctl restart <your-service>
#
# All child processes of the service will also inherit LD_PRELOAD,
# so a bash launcher script that spawns multiple processes is fully covered.

[Service]
Environment=LD_PRELOAD=/usr/lib/libcrashomon.so
Environment=CRASHOMON_DB_PATH=/var/crashomon
Environment=CRASHOMON_HANDLER_PATH=/usr/libexec/crashomon/crashpad_handler
